# ops/systemd/supplycore-lane-realtime.service
[Unit]
Description=SupplyCore lane: realtime (latency-sensitive syncs, dashboards, alerts)
After=network.target mariadb.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/SupplyCore
EnvironmentFile=-/etc/default/supplycore-worker
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/SupplyCore/.venv-orchestrator/bin/python /var/www/SupplyCore/bin/python_orchestrator.py loop-runner --app-root /var/www/SupplyCore --lane realtime --max-parallel 4 --fast-pause 5 --background-pause 15
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=90
MemoryMax=6G
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
