# ops/systemd/supplycore-lane-compute-bg.service
[Unit]
Description=SupplyCore lane: compute background loop (backfills, forecasting, CIP pipeline)
After=network.target mariadb.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/SupplyCore
EnvironmentFile=-/etc/default/supplycore-worker
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/SupplyCore/.venv-orchestrator/bin/python /var/www/SupplyCore/bin/python_orchestrator.py loop-runner --app-root /var/www/SupplyCore --lane compute --max-parallel 2 --background-pause 30 --background-only --memory-max-gb 24.0 --no-tier-barriers
Nice=10
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=90
MemoryMax=32G
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
