# ops/systemd/supplycore-loop-runner.service
[Unit]
Description=SupplyCore loop runner (all lanes, monolithic fallback)
After=network.target mariadb.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/SupplyCore
EnvironmentFile=-/etc/default/supplycore-worker
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/SupplyCore/.venv-orchestrator/bin/python /var/www/SupplyCore/bin/python_orchestrator.py loop-runner --app-root /var/www/SupplyCore --max-parallel 6
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=90
MemoryMax=4G
StandardOutput=append:/var/www/SupplyCore/storage/logs/worker.log
StandardError=append:/var/www/SupplyCore/storage/logs/worker.log

[Install]
WantedBy=multi-user.target
