# ops/systemd/supplycore-zkill.service
[Unit]
Description=SupplyCore zKill stream worker
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/SupplyCore
Environment=SUPPLYCORE_ZKILL_POLL_SLEEP=10
EnvironmentFile=-/etc/default/supplycore-worker
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/SupplyCore/.venv-orchestrator/bin/python /var/www/SupplyCore/bin/python_orchestrator.py zkill-worker --app-root /var/www/SupplyCore --poll-sleep ${SUPPLYCORE_ZKILL_POLL_SLEEP}
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=90
MemoryMax=768M
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
