# ops/systemd/supplycore-influx-rollup-export.service
[Unit]
Description=SupplyCore InfluxDB historical rollup export
After=network-online.target mariadb.service influxdb.service
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/SupplyCore
EnvironmentFile=-/etc/default/supplycore-influx-rollup-export
ExecStart=/var/www/SupplyCore/.venv-orchestrator/bin/python /var/www/SupplyCore/bin/python_orchestrator.py influx-rollup-export --app-root /var/www/SupplyCore $SUPPLYCORE_INFLUX_EXPORT_ARGS
Nice=10
IOSchedulingClass=best-effort
StandardOutput=append:/var/www/SupplyCore/storage/logs/influx-rollup-export.log
StandardError=append:/var/www/SupplyCore/storage/logs/influx-rollup-export.log

[Install]
WantedBy=multi-user.target
