# ops/systemd/supplycore-loop-runner.service
[Unit]
Description=SupplyCore loop runner (tier-by-tier job execution)
After=network.target mariadb.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/SupplyCore
EnvironmentFile=-/etc/default/supplycore-worker
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/SupplyCore/.venv-orchestrator/bin/python /var/www/SupplyCore/bin/python_orchestrator.py loop-runner --app-root /var/www/SupplyCore --max-parallel 2 --fast-pause 5 --background-pause 30
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=120
MemoryMax=5G
StandardOutput=append:/var/www/SupplyCore/storage/logs/worker.log
StandardError=append:/var/www/SupplyCore/storage/logs/worker.log

# NOTE: Disabled in favour of supplycore-sync-worker + supplycore-compute-worker.
# The loop-runner runs all ~64 jobs in a single process with thread pools,
# which causes OOM kills even at MemoryMax=5G.  The worker-pool model runs
# one job at a time per process with proper memory isolation.
[Install]
# WantedBy=multi-user.target
